# etc/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=false
autorestart=false
startsecs=0
